# theme/luci-theme-advanced-material/root/etc/uci-defaults/30_luci-theme-advanced-material
#!/bin/sh
uci batch <<-EOF
	set luci.themes.advanced-material=/luci-static/advanced-material
	set luci.main.mediaurlbase=/luci-static/advanced-material
	commit luci
EOF
exit 0
